# supervisor/kubernetes/test/conf.d/fastapi.conf
[program:fastapi-102_replace_this]
command=python -m debugpy --listen 10001 -m uvicorn --reload --host 0.0.0.0 --port 102_replace_this --log-level debug main:app
user=fastapi-skeleton
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
directory=/app
autostart=true
autorestart=true
stopsignal=QUIT
environment=ENVIRONMENT=%(ENV_ENVIRONMENT)s,ENVNAME=%(ENV_ENVNAME)s
